# crontab
15,45 * * * *  root python /srv/crawler.py /srv/ >> /var/log/cron.log 2>&1
